# duckdns.service
[Unit]
Description=DuckDNS update job
After=network.target

[Service]
Type=oneshot
ExecStart=/usr/bin/duckdns
