# duckdns.timer
[Unit]
Description=Run DuckDNS every 5 minutes

[Timer]
OnBootSec=1
OnUnitActiveSec=5min

[Install]
WantedBy=timers.target
